# update-wn-clients.cron
RANDOM_DELAY = 10
# ^ minutes
0 */12 * * * root /usr/bin/update-all-remote-wn-clients

